# rc.d/mdns_repeater
#!/bin/ksh

# Example OpenBSD rc.d(8) service for mdns-repeater.
daemon="/usr/local/sbin/mdns-repeater"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
